# versebg.service
[Unit]
Description=Verse of the day background generator

[Service]
ExecStart=/usr/bin/versebgd

[Install]
WantedBy=multi-user.target